# scripts/agsml_oracle.service
[Unit]
Description=agsml service

[Service]
ExecStart=/home/opc/agsml/scripts/agsml_server_start.sh

[Install]
WantedBy=multi-user.target